# scripts/com.baymax.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <!-- Unique identifier for this Launch Agent -->
  <key>Label</key>
  <string>com.baymax.daemon</string>

  <!-- The command to run: python daemon.py from the project root -->
  <key>ProgramArguments</key>
  <array>
    <string>PYTHON_PATH_PLACEHOLDER</string>
    <string>BAYMAX_DIR_PLACEHOLDER/daemon.py</string>
  </array>

  <!-- Working directory = project root -->
  <key>WorkingDirectory</key>
  <string>BAYMAX_DIR_PLACEHOLDER</string>

  <!-- Start automatically on login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart automatically if it crashes -->
  <key>KeepAlive</key>
  <true/>

  <!-- Wait 10 seconds before restarting after crash -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <!-- Log stdout and stderr to file -->
  <key>StandardOutPath</key>
  <string>HOME_PLACEHOLDER/Library/Logs/baymax.log</string>
  <key>StandardErrorPath</key>
  <string>HOME_PLACEHOLDER/Library/Logs/baymax_error.log</string>

  <!-- Environment variables (augment PATH so Python is found) -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>HOME_PLACEHOLDER</string>
  </dict>

</dict>
</plist>
